# exporter/etc/cron.d/exporter
7 15 * * * root php /code/exporter/db-diff/bin/pdo-git export      MarketflowAcc   Marketflow..TITRE  TIT_COD
9 15 * * * root php /code/exporter/db-diff/bin/pdo-git export      MarketflowAccDB MkfFFADubai..TITRE TIT_COD
11 15 * * * root php /code/exporter/db-diff/bin/pdo-git post-commit MarketflowAcc   Marketflow..TITRE  --columns /code/exporter/db-diff/src/columns/ffa-titre.yml --format email --email.to s.akiki@ffaprivatebank.com --email.config /code/exporter/etc/swiftmailer/ffa.yml
12 15 * * * root php /code/exporter/db-diff/bin/pdo-git post-commit MarketflowAccDB MkfFFADubai..TITRE --columns /code/exporter/db-diff/src/columns/ffa-titre.yml --format email --email.to s.akiki@ffaprivatebank.com --email.config /code/exporter/etc/swiftmailer/ffa.yml
